# meta-ami/recipes-phosphor/gpio/crashdump-trigger/crashdump-trigger.service
[Unit]
Description=crashdump Trigger Service

[Service]
ExecStart=/usr/bin/crashdump_trigger.sh

[Install]
WantedBy=multi-user.target
