# openbmc_modules/x86-power-control/power-control-x86/service_files/chassis-system-reset.target
[Unit]
Description=System unit to hard reset or system reset chassis
