# openbmc_modules/x86-power-control/power-control-x86/service_files/xyz.openbmc_project.Chassis.Control.Power.service
[Unit]
Description=Intel Power Control

[Service]
Restart=always
RestartSec=3
ExecStart=/usr/bin/power-control
Type=dbus
BusName=xyz.openbmc_project.State.Host

[Install]
WantedBy=sysinit.target

