# meta-ami/meta-aowanda/recipes-phosphor/gpio/host-powerstatus-trigger/host_status-trigger.service
[Unit]
Description=Host power status trigger service

[Service]
ExecStart=/usr/bin/host_status_trigger.sh

[Install]
WantedBy=multi-user.target

