# meta-ami/meta-aowanda/recipes-phosphor/gpio/host-powerstatus-trigger/hostpower-gpio-monitor.service
  
[Unit]
Description=Phosphor Multi GPIO monitor

[Service]
ExecStart=/usr/bin/phosphor-multi-gpio-monitor --config /usr/share/phosphor-gpio-monitor/host_status_trigger.json

[Install]
WantedBy=multi-user.target

