# meta-ami/meta-common/recipes-ami/bmc-boot-check/files/xyz.openbmc_project.bmcbootcheck.service
[Unit]
Description=Check if bmc reboot is caused by AC loss

[Service]
RemainAfterExit=true
Type=oneshot
ExecStart=/usr/bin/env bmc-boot-check.sh

[Install]
WantedBy=multi-user.target
