# openbmc_modules/phosphor-net-ipmid/phosphor-ipmi-net@.socket
[Socket]
ListenDatagram=623
BindToDevice=sys-subsystem-net-devices-%i.device

[Install]
WantedBy=sockets.target

