# openbmc_modules/phosphor-bmc-code-mgmt/mmc/obmc-flash-mmc-remove@.service.in
[Unit]
Description=Delete image %I from BMC storage

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc mmc-remove %i
