# openbmc_modules/phosphor-bmc-code-mgmt/mmc/obmc-flash-mmc-setprimary@.service.in
[Unit]
Description=Set %I as primary version

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc mmc-setprimary %i
