# openbmc_modules/phosphor-ipmi-flash/bmc/phosphor-ipmi-flash-bios-prepare.target.in
[Unit]
Description=Phosphor-ipmi-flash Prepare BIOS to receive update
