# openbmc_modules/phosphor-ipmi-flash/bmc/phosphor-ipmi-flash-bmc-update.target.in
[Unit]
Description=Phosphor-ipmi-flash update the BMC image
