# openbmc_modules/phosphor-ipmi-flash/bmc/phosphor-ipmi-flash-bios-update.target.in
[Unit]
Description=Phosphor-ipmi-flash update the BIOS image
