# openbmc_modules/phosphor-ipmi-flash/bmc/phosphor-ipmi-flash-bmc-verify.target.in
[Unit]
Description=Phosphor-ipmi-flash verify the image contents
