# openbmc_modules/phosphor-ipmi-flash/bmc/phosphor-ipmi-flash-bmc-prepare.target.in
[Unit]
Description=Phosphor-ipmi-flash Prepare BMC to receive update
