# meta-ami/recipes-phosphor/intelacd/files/IntelACD.service
[Unit]
Description= Intel ACD start script

[Service]
ExecStart=/usr/bin/crashdump

[Install]
WantedBy=multi-user.target
