# openbmc_modules/phosphor-state-manager/target_files/obmc-power-start-pre@.target
[Unit]
Description=Power%i On (Pre)
After=multi-user.target
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
