# openbmc_modules/phosphor-state-manager/target_files/obmc-host-quiesce@.target
[Unit]
Description=Quiesce Target
After=multi-user.target
RefuseManualStop=yes
Conflicts=obmc-chassis-poweroff@%i.target
