# openbmc_modules/phosphor-state-manager/target_files/obmc-host-start-pre@.target
[Unit]
Description=Start Host%i (Pre)
After=multi-user.target
Conflicts=obmc-host-stop@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
