# openbmc_modules/phosphor-state-manager/target_files/obmc-host-startmin@.target
[Unit]
Description=Start Host%i Minimum
After=multi-user.target
Conflicts=obmc-host-stop@%i.target
RefuseManualStop=yes
OnFailure=obmc-host-quiesce@%i.target
OnFailureJobMode=flush
