# openbmc_modules/phosphor-state-manager/target_files/obmc-host-started@.target
[Unit]
Description=Host%i (Started)
After=multi-user.target
After=obmc-host-starting@%i.target
Conflicts=obmc-host-stop@%i.target
RefuseManualStart=yes
RefuseManualStop=yes
