# openbmc_modules/phosphor-bmc-code-mgmt/ubi/obmc-flash-bmc-ubirw-remove.service.in
[Unit]
Description=Deletes BMC read-write ubi volume

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc ubiremove rwfs
