# meta-ami/meta-jade/recipes-phosphor/gpio/scp-failover/scp_failover.service
[Unit]
Description=Ampere Computing LLC Ampere failover event service

[Service]
Restart=no
ExecStart=/usr/bin/scp_failover.sh 
