# openbmc_modules/phosphor-bmc-code-mgmt/obmc-flash-host-bios@.service.in
[Unit]
Description=Flash Host Bios image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=echo Please add custom command for flashing image /tmp/image/%i
