# openbmc_modules/phosphor-bmc-code-mgmt/obmc-flash-bmc-setenv@.service.in
[Unit]
Description=Set U-Boot environment variable

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc ubisetenv %I
