# openbmc_modules/phosphor-bmc-code-mgmt/reboot-guard-disable.service.in
[Unit]
Description=Removes the guard that blocks BMC reboot

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc rebootguarddisable
