# openbmc_modules/phosphor-bmc-code-mgmt/xyz.openbmc_project.Software.Download.service.in
[Unit]
Description=Phosphor Download Manager

[Service]
ExecStart=/usr/bin/phosphor-download-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
