# openbmc_modules/obmc-console/conf/obmc-console-ssh.socket.in
[Unit]
Description=Phosphor Host Console SSH Per-Connection socket

[Socket]
ListenStream=2200
Accept=yes

[Install]
WantedBy=sockets.target
