# meta-ami/recipes-phosphor/network/set-hostname/ntp-restart.path
[Path]
Unit=ntp-restart.service
PathModified=/etc/systemd/network/00-bmc-eth0.network

[Install]
WantedBy=multi-user.target
