# meta-ami/recipes-phosphor/network/set-hostname/register-dns.path
[Path]
Unit=register-dns.service
PathModified=/etc/hostname

[Install]
WantedBy=multi-user.target
