# openbmc_modules/ipmbbridge/ipmb.service
[Unit]
Description=IPMB bridge
After=phosphor-ipmi-host.service

[Service]
ExecStart=/usr/bin/ipmbbridged
SyslogIdentifier=ipmbbridged
Restart=always
Type=simple

[Install]
WantedBy=basic.target
